# gunicorn-www.control-net.org.conf
start on net-device-up
stop on shutdown

respawn
setuid caffe
chdir /home/caffe/mysite
exec ../env/bin/gunicorn --bind unix:/tmp/www.control-net.org.socket mysite.wsgi:application

